# deploy/mdvwb-manager.service
[Unit]
Description=MDVWB multi-bus configuration and MQTT manager
Wants=network-online.target mosquitto.service
After=network-online.target mosquitto.service
Before=mdvwb@1.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/mdvwb-manager
ExecStart=/usr/local/bin/mdvwb-manager mqtt /etc/mdvwb/buses.json
Restart=on-failure
RestartSec=5
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
